# Configs/.config/autostart/OneDrive_Mount.sh.desktop
[Desktop Entry]
Comment[en_US]=
Comment=
Exec=~/.local/share/OneDrive_Mount.sh
GenericName[en_US]=
GenericName=
Icon=application-x-shellscript
MimeType=
Name[en_US]=OneDrive_Mount.sh
Name=OneDrive_Mount.sh
Path=
PrefersNonDefaultGPU=false
StartupNotify=true
Terminal=false
TerminalOptions=
Type=Application
X-KDE-AutostartScript=true
X-KDE-SubstituteUID=false
X-KDE-Username=
